# templates/rocket.service.erb
[Unit]
Description=Rocket.Chat service

[Service]
ExecStart=/usr/local/bin/node <%= @destination %>/bundle/main.js
Restart=always
# Restart service after 10 seconds if node service crashes
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rocket.chat
Environment=PORT=<%= @port %> MONGO_URL=mongodb://<%= @mongo_host %>:<%= @mongo_port%>/<%= @database_name%> ROOT_URL=<%= @root_url %>

[Install]
WantedBy=multi-user.target
