# templates/rocket.service.erb
<%- # guards
authstring = @authsource.to_s != '' ? "?authSource=#{@authsource}" : nil
replstring = @mongo_replset.to_s != '' ? (authstring ? '&' : '?') + "replicaSet=#{@mongo_replset}" : nil
-%>
[Unit]
Description=Rocket.Chat service

[Service]
ExecStart=/usr/local/bin/node <%= @destination %>/bundle/main.js
Restart=always
# Restart service after 10 seconds if node service crashes
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=rocket.chat
Environment=INSTANCE_IP=<%= @instance_ip || scope.lookupvar('networking::ip') %>
Environment=PORT=<%= @port %>
Environment=MONGO_URL=mongodb://<%= @mongo_host%>:<%= @mongo_port%>/<%= @database_name%><%= authstring %><%= replstring %>
Environment=ROOT_URL=<%= @root_url %>
<%- if @mongo_replset -%>
Environment=MONGO_OPLOG_URL=mongodb://<%= @mongo_host%>:<%= @mongo_port%>/local<%= authstring %><%= replstring %>
<%- end -%>

[Install]
WantedBy=multi-user.target
